# ts-vanilla/etc/supervisord/service-name.conf
[supervisord]
nodaemon=true

[program:nginx]
command=service nginx start
stopsignal=KILL
autostart=true
autorestart=true
startsecs=10
